# docker/supervisord.conf
[supervisord]
nodaemon = true

[program:postgres]
user = postgres
command = /usr/lib/postgresql/9.4/bin/postgres -D /var/lib/postgresql/9.4/main -c config_file=/etc/postgresql/9.4/main/postgresql.conf

[program:uwsgi]
command = /s/bin/uwsgi --die-on-term --emperor /srv/uwsgi/conf/uwsgi

